# rc.local
#!/bin/sh -e
#
# rc.local
#

dhclient eth0;
ifconfig eth0 192.168.1.70/24;
ifconfig eth0:2 10.10.10.124;

sleep 3;

bash -c "screen -d -m /dev/ttyUSB0 &";

sleep 10;

killall -9 screen;

bash -c "motion";


echo -n "wwwwww" > /dev/ttyUSB0;




exit 0
